# hosts/linux/omarchy/files/99-ds740d-force-usb2.rules
# Brother DS-740D (USB 04f9:0469) is unstable/undriveable at USB 3 SuperSpeed:
# spontaneous disconnects, and the brscan5 image read fails with "Error during
# device I/O". Force it to USB 2.0 by disabling the SuperSpeed side of its root
# port; it then falls back to the High-Speed peer and scans reliably.
#
# The scanner sits on a fixed internal hub port: SuperSpeed instance = 2-1.2 on
# root port usb2-port1; High-Speed peer = 1-1.2. (Paths are topology-specific —
# if the scanner is moved to another port these must be updated.)
#
# Two rules, race-free:
#  - on SuperSpeed connect (speed 5000): disable SS → it re-enumerates at HS.
#  - on HIGH-SPEED disconnect (kernel 1-1.2 = a real unplug/sleep, NOT the
#    bump-induced SS disconnect): re-enable SS so the next connect can bump again.
#    Matching only the HS instance is what avoids the earlier race that bounced
#    the scanner off the bus.
#
# Root-owned (home-manager can't manage /etc on a foreign distro); see
# hosts/linux/omarchy/default.nix for the one-time install.
ACTION=="add", SUBSYSTEM=="usb", ATTR{idVendor}=="04f9", ATTR{idProduct}=="0469", ATTR{speed}=="5000", RUN+="/bin/sh -c 'echo 1 > /sys/bus/usb/devices/usb2/2-0:1.0/usb2-port1/disable'"
ACTION=="remove", SUBSYSTEM=="usb", KERNEL=="1-1.2", RUN+="/bin/sh -c 'echo 0 > /sys/bus/usb/devices/usb2/2-0:1.0/usb2-port1/disable'"
