# hosts/linux/omarchy/files/ds740d-force-usb2.service
[Unit]
# Brother DS-740D (USB 04f9:0469) is unstable/undriveable at USB 3 SuperSpeed:
# it spontaneously disconnects and the brscan5 image read fails with "Error
# during device I/O". Force its port to USB 2.0 by disabling the SuperSpeed side
# of its root port (usb2-port1) — the device then falls back to the High-Speed
# peer (usb1-port1) and scans reliably. Only the scanner is on this port.
#
# Root-owned system unit (installed out-of-band; home-manager can't manage /etc
# on a foreign distro). See hosts/linux/omarchy/default.nix for the install step.
Description=Force Brother DS-740D USB port to USB 2.0
After=multi-user.target

[Service]
Type=oneshot
RemainAfterExit=yes
# The scanner is bus-powered and enumerates at boot; write the disable after USB
# settles. Harmless if the port path is absent (returns success either way).
ExecStart=/bin/sh -c 'echo 1 > /sys/bus/usb/devices/usb2/2-0:1.0/usb2-port1/disable 2>/dev/null || true'

[Install]
WantedBy=multi-user.target
